# alert-bridge/trade-poller.service
[Unit]
Description=Project Nate - DPL Alert Poller (monitors Nate Bear's trades)
After=network.target trade-webhook.service
Requires=trade-webhook.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/nate.tps.pro/alert-bridge
ExecStart=/root/nate.tps.pro/.venv/bin/python /root/nate.tps.pro/alert-bridge/alert_poller.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
