# alert-bridge/trade-poller.service
[Unit]
Description=DailyProfitsLive Alert Poller
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/nate.tps.pro
EnvironmentFile=/root/nate.tps.pro/.env.trade
ExecStart=/root/nate.tps.pro/.venv/bin/python /root/nate.tps.pro/alert-bridge/alert_poller.py
Restart=always
RestartSec=10
StandardOutput=append:/root/nate.tps.pro/logs/poller.log
StandardError=append:/root/nate.tps.pro/logs/poller.log

[Install]
WantedBy=multi-user.target
