# alert-bridge/trade-webhook.service
[Unit]
Description=Project Nate - Webhook Server (receives alerts from poller)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/nate.tps.pro/alert-bridge
ExecStart=/root/nate.tps.pro/.venv/bin/python /root/nate.tps.pro/alert-bridge/webhook_server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
